# rpisolark-monitor.service
[Unit]
Description=RpiSolarkMonitor - Production frequency monitor
After=network.target

[Service]
Type=simple
# Run as root (no sudo needed) or as your user (requires passwordless sudo)
# Option 1: Run as root (simplest, no sudo needed in script)
User=root
WorkingDirectory=/home/keithcu/Desktop/RpiSolArk
# Explanation:
# - We will call a wrapper script that:
#   - Resolves its own directory dynamically
#   - Uses that as the app/.venv directory
#   - Activates .venv if present (or runs system python as fallback)
#   - Executes: sudo python monitor.py --real
ExecStart=/home/keithcu/Desktop/RpiSolArk/rpisolark-monitor.sh

Restart=on-failure
RestartSec=10
WatchdogSec=60
StartLimitIntervalSec=300
StartLimitBurst=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target