# rpisolark-monitor.service
[Unit]
Description=RpiSolarkMonitor - Production frequency monitor
After=network.target

[Service]
Type=simple
# Run as 'pi' by default; adjust if your deployment user differs.
User=pi
# Use the directory containing this service's companion script as WorkingDirectory.
# This assumes you place rpisolark-monitor.sh and this unit in the same install directory.
WorkingDirectory=%E
# Explanation:
# - We will call a wrapper script that:
#   - Resolves its own directory dynamically
#   - Uses that as the app/.venv directory
#   - Activates .venv if present (or runs system python as fallback)
#   - Executes: python monitor.py --real
ExecStart=%E/rpisolark-monitor.sh

Restart=on-failure
RestartSec=10
WatchdogSec=60
StartLimitIntervalSec=300
StartLimitBurst=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target