# update.timer
[Unit]
Description=Run auto-update every hour
Requires=update.service

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
